# agent-cluster/deployment/agent-cluster.service
#!/bin/bash
# Agent Cluster V2.7.1 开机自启脚本
# 部署到：/etc/systemd/system/agent-cluster.service

[Unit]
Description=Agent Cluster V2.7.1 Web Service
Documentation=https://github.com/phoenixbull/agent-cluster-test
After=network.target openclaw.service
Wants=network.target

[Service]
Type=simple
User=admin
Group=admin
WorkingDirectory=/home/admin/.openclaw/workspace/agent-cluster
Environment="PATH=/usr/bin:/bin:/usr/local/bin:/home/admin/.local/bin"
Environment="PYTHONUNBUFFERED=1"

# 启动主 Web 服务
ExecStart=/usr/bin/python3 /home/admin/.openclaw/workspace/agent-cluster/web_app_v2.py --port 8890
Restart=always
RestartSec=10
StartLimitInterval=60s
StartLimitBurst=3

# 启动钉钉通知服务
ExecStartPost=/bin/bash -c 'sleep 5 && /usr/bin/python3 /home/admin/.openclaw/workspace/agent-cluster/monitoring/dingtalk_notifier.py &'

# 日志
StandardOutput=journal
StandardError=journal
SyslogIdentifier=agent-cluster

# 安全加固
NoNewPrivileges=true
PrivateTmp=true

# 资源限制
LimitNOFILE=65535
LimitNPROC=4096

# 工作目录权限
ReadWritePaths=/home/admin/.openclaw/workspace/agent-cluster

[Install]
WantedBy=multi-user.target
